# deployment/peerai-unified.service
[Unit]
Description=Peer AI Unified Service
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/peer-ai
ExecStart=/bin/bash /home/ubuntu/peer-ai/deployment/start_single_server.sh
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal
SyslogIdentifier=peerai-unified

# Security settings
PrivateTmp=true
ProtectSystem=full
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target 